# crontab

*/30 * * * * python app/producer.py >> /logfile